# rpi_ble_logger_service/ble_logger.service
[Unit]
Description=BLE logger

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/bin/bt_getter.py

[Install]
WantedBy=multi-user.target